# docker-simple/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true

[program:backend]
command=/app/airboard-api
directory=/app
autorestart=true
environment=GIN_MODE=release